# assets/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
username=nobody
password=nobody

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=1MB
logfile_backups=3
loglevel=error
pidfile=/tmp/supervisord.pid
nodaemon=true
minfds=10000
minprocs=2000

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:php-fpm]
command=/usr/sbin/php-fpm7.4 --fpm-config=/app/assets/php-fpm.conf
autostart=true
autorestart=true
startretries=5
priority=5
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true
startretries=5
priority=10
stopsignal=QUIT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[eventlistener:supervisord-helper]
command=/app/assets/supervisord-helper.py
events=PROCESS_STATE_FATAL
autostart=true
priority=15
stopsignal=QUIT
